# config/systemd/user/sunshine.service
[Unit]
Description=Sunshine game streaming

[Service]
Type=simple
StandardOutput=journal
ExecStart=sunshine

[Install]
WantedBy=default.target
